# camerad.service
[Unit]
Description=AutomaticCamera

[Service]
ExecStart=/opt/DCIM/camerad.sh

[Install]
WantedBy=multi-user.target
